# Bspwm: tea/.config/autostart/Picom.desktop
[Desktop Entry]
Exec=/usr/bin/picom
Name=Picom
Type=Application
Version=1.0
